# deploy/systemd/luckyred-bot.service
[Unit]
Description=LuckyRed Telegram Bot
After=network.target postgresql.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/luckyred/bot
Environment="PATH=/opt/luckyred/bot/.venv/bin"
EnvironmentFile=/opt/luckyred/.env
ExecStart=/opt/luckyred/bot/.venv/bin/python main.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target

