# deploy/systemd/luckyred-bot.service
[Unit]
Description=LuckyRed Telegram Bot
After=network.target postgresql.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/home/ubuntu/hbgm001/bot
Environment="PATH=/home/ubuntu/hbgm001/bot/.venv/bin"
EnvironmentFile=/home/ubuntu/hbgm001/.env
ExecStart=/home/ubuntu/hbgm001/bot/.venv/bin/python main.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target

