# deploy/systemd/luckyred-api.service
[Unit]
Description=LuckyRed MiniApp API Service
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/home/ubuntu/hbgm001/api
Environment="PATH=/home/ubuntu/hbgm001/api/.venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/home/ubuntu/hbgm001"
EnvironmentFile=/home/ubuntu/hbgm001/.env

# 启动命令
ExecStart=/home/ubuntu/hbgm001/api/.venv/bin/uvicorn main:app --host 127.0.0.1 --port 8080 --workers 4

# 重启策略
Restart=always
RestartSec=5
StartLimitInterval=200
StartLimitBurst=5

# 日志配置
StandardOutput=journal
StandardError=journal
SyslogIdentifier=luckyred-api

# 安全配置
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/ubuntu/hbgm001

[Install]
WantedBy=multi-user.target
