# deploy/systemd/luckyred-bot.service
[Unit]
Description=LuckyRed Telegram Bot
After=network.target postgresql.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/hbgm001
Environment="PATH=/home/ubuntu/hbgm001/bot/.venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/home/ubuntu/hbgm001"
EnvironmentFile=/home/ubuntu/hbgm001/.env
ExecStart=/home/ubuntu/hbgm001/bot/.venv/bin/python bot/main.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

