# deploy/systemd/luckyred-admin.service
[Unit]
Description=LuckyRed Web Admin
After=network.target postgresql.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/home/ubuntu/hbgm001/admin
Environment="PATH=/home/ubuntu/hbgm001/admin/.venv/bin"
EnvironmentFile=/home/ubuntu/hbgm001/.env
ExecStart=/home/ubuntu/hbgm001/admin/.venv/bin/uvicorn main:app --host 127.0.0.1 --port 8081
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target

